# DEBIAN/postinst
#!/bin/bash
set -euo pipefail

NC='\033[0m' # Sin color
AMARILLO='\033[1;33m'

# Recargar systemd y desplegar el servicio si systemctl está disponible
if command -v systemctl >/dev/null 2>&1; then
    echo "Recargando systemd y habilitando servicio..."
    systemctl daemon-reload || true
    systemctl unmask ssp.service || true
    if systemctl enable --now ssp.service; then
        echo "ssp.service habilitado y arrancado."
    else
        echo -e "${AMARILLO}ADVERTENCIA${NC}: No se pudo habilitar o arrancar ssp.service. Comprueba los logs de systemd."
    fi
else
    echo -e "${AMARILLO}ADVERTENCIA${NC}: systemctl no disponible: omitiendo habilitación/arranque del servicio."
fi

# Aplicar permisos de forma condicional (no fallar si los archivos no existen)
if [ -f "/lib/systemd/system/ssp.service" ]; then
    chmod 0644 /lib/systemd/system/ssp.service
else
    echo "Aviso: /lib/systemd/system/ssp.service no encontrado."
fi

if [ -f "/usr/sbin/ssp" ]; then
    chmod 0755 /usr/sbin/ssp
else
    echo -e "${AMARILLO}ADVERTENCIA${NC}: /usr/sbin/ssp no encontrado."
fi

if [ -f "/usr/sbin/ssp_files/service.py" ]; then
    chmod 0755 /usr/sbin/ssp_files/service.py
else
    echo -e "${AMARILLO}ADVERTENCIA${NC}: /usr/sbin/ssp_files/service.py no encontrado."
fi

echo "Postinst terminado"
exit 0