# DEBIAN/postinst
#!/bin/bash
set -euo pipefail

NC='\033[0m' # Sin color
AMARILLO='\033[1;33m'

# Recargar systemd y desplegar el servicio si systemctl está disponible
if command -v systemctl >/dev/null 2>&1; then
    echo "Recargando systemd y habilitando servicio..."
    systemctl daemon-reload || true
    systemctl unmask victrix.service || true
    if systemctl enable --now victrix.service; then
        echo "victrix.service habilitado y arrancado."
    else
        echo -e "${AMARILLO}ADVERTENCIA${NC}: No se pudo habilitar o arrancar victrix.service. Comprueba los logs de systemd."
    fi
else
    echo -e "${AMARILLO}ADVERTENCIA${NC}: systemctl no disponible: omitiendo habilitación/arranque del servicio."
fi

# Aplicar permisos de forma condicional (no fallar si los archivos no existen)
if [ -f "/lib/systemd/system/victrix.service" ]; then
    chmod 0644 /lib/systemd/system/victrix.service
else
    echo "Aviso: /lib/systemd/system/victrix.service no encontrado."
fi

if [ -f "/usr/sbin/victrix" ]; then
    chmod 0755 /usr/sbin/victrix
else
    echo -e "${AMARILLO}ADVERTENCIA${NC}: Fichero /usr/sbin/victrix no encontrado."
fi

if [ -d "/usr/sbin/victrix_files" ]; then
    chmod 0755 /usr/sbin/victrix_files
else
    echo -e "${AMARILLO}ADVERTENCIA${NC}: Directorio /usr/sbin/victrix_files no encontrado."
fi

echo "Postinst terminado"
exit 0