# DEBIAN/postinst
#!/bin/bash

whitelist="/usr/sbin/ssp_files/deb_services"

set -euo pipefail

# Limpieza de consola (opcional)
# clear

echo "Iniciando postinst"

# Crear directorio para archivos del servicio si no existe
install -d -m 0755 /usr/sbin/ssp_files

# Instalar la licencia en /usr/sbin/ssp_files como LICENCIA_Apache2.0
# Intentar varias rutas relativas/ya instaladas para mayor robustez
if [ -f "../LICENSE" ]; then
    install -m 0644 ../LICENSE /usr/sbin/ssp_files/LICENCIA_Apache2.0
    echo "License installed from ../LICENSE"
elif [ -f "/usr/sbin/ssp_files/LICENSE" ]; then
    install -m 0644 /usr/sbin/ssp_files/LICENSE /usr/sbin/ssp_files/LICENCIA_Apache2.0
    echo "License installed from /usr/sbin/ssp_files/LICENSE"
elif [ -f "./LICENSE" ]; then
    install -m 0644 ./LICENSE /usr/sbin/ssp_files/LICENCIA_Apache2.0
    echo "License installed from ./LICENSE"
else
    echo "WARNING: LICENSE not found. Skipping license installation."
fi

# Preguntar al usuario si desea agregar servicios de usuario (altamente recomendado)
read -p "Usuario '$user' detectado, ¿desea agregarlo a la lista blanca? (recomendado) [s/n]: " response
response=$(echo "$response" | tr '[:upper:]' '[:lower:]')
if [[ "$response" == "s" || "$response" == "si" ]]; then
    echo "Agregando servicios del usuario '$user' a la lista blanca..."
    # Agregar servicios user@ y user-runtime-dir@
    systemctl list-units --type=service --state=active --no-legend \
      | awk '{print $1}' \
      | grep -E '^user(@|-runtime-dir@)[0-9]+\.service$' \
      >> "$whitelist"
    echo "Servicios agregados correctamente a: $whitelist"
else
    echo "[PELIGRO] Podría tener problemas para iniciar sesión si no agrega los servicios de usuario."
    read -p "¿Está seguro de que NO desea agregarlos a la lista blanca? [s/n]: " confirm_response
    confirm_response=$(echo "$confirm_response" | tr '[:upper:]' '[:lower:]')
    if [[ "$confirm_response" == "n" || "$confirm_response" == "no" ]]; then
        echo "Agregando servicios del usuario '$user' a la lista blanca (por seguridad)..."

        systemctl list-units --type=service --state=active --no-legend \
          | awk '{print $1}' \
          | grep -E '^user(@|-runtime-dir@)[0-9]+\.service$' \
          >> "$whitelist"

        echo "Servicios agregados correctamente a: $whitelist"
    else
        echo "No se agregaron servicios. Procediendo bajo su responsabilidad..."
    fi
fi


# Recargar systemd y desplegar el servicio
echo "Recargando systemd y habilitando servicio..."
# daemon-reload puede fallar si no hay privilegios; en postinst normalmente corre como root
systemctl daemon-reload || true
systemctl unmask ssp.service || true
# Intentar habilitar y arrancar el servicio inmediatamente
if systemctl enable --now ssp.service; then
    echo "ssp.service habilitado y arrancado."
else
    echo "WARNING: No se pudo habilitar o arrancar ssp.service. Comprueba los logs de systemd."
fi

echo "Aplicando permisos..."
# Aplicar permisos de forma condicional (no fallar si los archivos no existen)
if [ -f "/lib/systemd/system/ssp.service" ]; then
    chmod 0644 /lib/systemd/system/ssp.service
else
    echo "Aviso: /lib/systemd/system/ssp.service no encontrado."
fi

if [ -f "/usr/sbin/ssp" ]; then
    chmod 0755 /usr/sbin/ssp
else
    echo "Aviso: /usr/sbin/ssp no encontrado."
fi

if [ -f "/usr/sbin/ssp_files/service.py" ]; then
    chmod 0755 /usr/sbin/ssp_files/service.py
else
    echo "Aviso: /usr/sbin/ssp_files/service.py no encontrado."
fi

echo "Postinst terminado"
exit 0