# DEBIAN/prerm
#!/bin/bash
set -euo pipefail

# prerm: se ejecuta antes de la eliminación del paquete
# Detener y deshabilitar servicio de forma segura si systemctl está disponible
if command -v systemctl >/dev/null 2>&1; then
    echo "Deteniendo y deshabilitando victrix.service si existe..."
    # Intentar parar (no fallar si no está activo)
    systemctl stop victrix.service >/dev/null 2>&1 || true
    # Intentar deshabilitar (no fallar si no está habilitado)
    systemctl disable victrix.service >/dev/null 2>&1 || true
else
    echo "systemctl no disponible: saltando parada/deshabilitado del servicio."
fi

# Eliminación segura de rutas del sistema (solo si existen)
echo "Eliminando archivos y directorios instalados por el paquete (si existen)..."

if [ -d "/usr/sbin/victrix_files" ]; then
    rm -rf "/usr/sbin/victrix_files"
    echo "Eliminado /usr/sbin/victrix_files"
fi

if [ -f "/usr/sbin/victrix" ]; then
    rm -f "/usr/sbin/victrix"
    echo "Eliminado /usr/sbin/victrix"
fi

# Ruta de manual corregida
if [ -f "/usr/share/man/man8/victrix.8" ]; then
    rm -f "/usr/share/man/man8/victrix.8"
    echo "Eliminado /usr/share/man/man8/victrix.8"
fi

if [ -d "/var/log/victrix" ]; then
    rm -rf "/var/log/victrix"
    echo "Eliminado /var/log/victrix"
fi

if [ -f "/lib/systemd/system/victrix.service" ]; then
    rm -f "/lib/systemd/system/victrix.service"
    echo "Eliminado /lib/systemd/system/victrix.service"
fi

if [ -f "/etc/apt/apt.conf.d/99update-victrix" ]; then
    rm -f "/etc/apt/apt.conf.d/99update-victrix"
    echo "Eliminado /etc/apt/apt.conf.d/99update-victrix"
fi

# Recargar systemd para que tenga en cuenta la eliminación de la unidad
if command -v systemctl >/dev/null 2>&1; then
    systemctl daemon-reload >/dev/null 2>&1 || true
fi

exit 0