# DEBIAN/prerm
#!/bin/bash

echo "Running prerm-install script..."

#systemctl stop ssp.service
#systemctl disable ssp.service